# systemd/ai-mkt-eval.timer
[Unit]
Description=Run AI-MKT eval nightly at 02:00

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
